# etc/systemd/system/snide.service
[Unit]
Description=SNI Demultiplexer
Requires=network.target
After=multi-user.target

[Service]
Type=notify
ExecStart=/usr/local/bin/snide -d /etc/snide

[Install]
WantedBy=multi-user.target
